# novicecpp/homeconfig/roles/firefox_desktop/files/dot-config/systemd/user/firefox_profiles_backup.service
[Unit]
Description=Backup firefox profile

[Service]
Type=simple
ExecStart=/opt/firefox_profiles_backup/backup.sh

[Install]
WantedBy=default.target
