# novicecpp/homeconfig/roles/firefox_desktop/files/dot-config/systemd/user/firefox_profiles_backup.timer
[Unit]
Description=Backup firefox profile every hour

[Timer]
OnCalendar=*-*-* 11:00:00
OnCalendar=*-*-* 23:00:00
Persistent=true

[Install]
WantedBy=timers.target
