# ansible_collections/novicecpp/homeconfig/roles/wireguard/files/etc/systemd/system/wg-quick@.service.d/override.conf
[Unit]
PartOf=ping-wg.service

[Service]
Restart=on-failure
RestartSec=10
