# novicecpp/homeconfig/roles/coredns/templates/etc/systemd/system/coredns.service.d/override.conf.j2
[Service]
ExecStart=
ExecStart=/opt/coredns/bin/coredns -conf=/etc/coredns/Corefile

[Unit]
After=wg-quick@{{ coredns_bind_interface }}.service
