# ansible_collections/novicecpp/homeserver/roles/certbot/files/etc/systemd/system/certbot.timer
[Unit]
Description=Run certbot weekly
After=network-online.target

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
