# novicecpp/homeconfig/roles/ssh_agent/files/dot-config/systemd/user/ssh-agent@.service
[Unit]
Description=ssh-agent for %I

[Service]
Type=simple
ExecStart=%h/.local/bin/ssh-agent.sh %i
Restart=on-failure
RestartSec=1s

[Install]
WantedBy=default.target
